# roles/docker/templates/docker.service.j2
# {{ansible_managed}}
[Service]
ExecStart=
ExecStart=/usr/bin/docker daemon -H fd:// --graph="/var/lib/{{ansible_hostanme}}/docker" --storage-driver=overlay -D
